# vcontrold/rootfs/etc/services.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare config_commands

config_commands=$(bashio::config 'commands')

## Fetch MQTT settings
. /etc/services.d/get_mqtt_settings.sh

## Fetch Vcontrold settings
. /etc/services.d/get_vcontrold_settings.sh

## Run your program
bashio::log.info "Starting vclient_pub main loop"
while sleep $REFRESH_RATE; do
    
#wether to run each comand separatley and open a new serial connection for each command or not
if bashio::config.true 'loop_vclient'; then 
    bashio::log.info "Starting one vclient per command"
    
    commands_array=(`echo $config_commands | tr '|' ' '`)
    i=0
    for command in "${commands_array[@]}"
    do
       :
        ((i=i+1))
        cmd=$(echo $command | cut -d':' -f1)
        type=$(echo $command | cut -d':' -f2)
        
        # Cleanup scripts
        rm /etc/vcontrold/1_mqtt_commands.txt || /bin/true
        rm /etc/vcontrold/2_mqtt.tmpl || /bin/true
        echo '#!/bin/sh' > /etc/vcontrold/2_mqtt.tmpl
        
        bashio::log.info "Run command ${i}: ${cmd} with type ${type}"
        echo $cmd >> /etc/vcontrold/1_mqtt_commands.txt
    
        R='$'$i''
        if [[ $type == "STRING" ]]; then
            R='"$R'$i'"'
        fi
    
        echo 'if [ "x$E'$i'" = "xOK" ]; then' >> /etc/vcontrold/2_mqtt.tmpl
        echo '  mosquitto_pub -h $MQTT_HOST -p $MQTT_PORT -u $MQTT_USER -P $MQTT_PASSWORD -t $MQTT_TOPIC/$C'$i' -m '$R'' >> /etc/vcontrold/2_mqtt.tmpl
        echo 'fi' >> /etc/vcontrold/2_mqtt.tmpl
   
        ## only one command per run, no list anymore
        cp -f /etc/vcontrold/2_mqtt.tmpl /config/vcontrold/2_mqtt.tmpl        
        vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -f /etc/vcontrold/1_mqtt_commands.txt -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh && exitcode=0 || exitcode=$?
        cp -f /etc/vcontrold/3_mqtt_pub.sh /config/vcontrold/3_mqtt_pub.sh
        bashio::log.info "vclient_pub line $cmd finished with exit code $exitcode."
        sleep 1
    done

else
    ## only one vclient per run, all commands are in the list 1_mqtt_commands.txt (faster, but maybe not so reliable)
    bashio::log.info "Starting one vclient for all commands"
    vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -f /etc/vcontrold/1_mqtt_commands.txt -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh && exitcode=0 || exitcode=$?
fi
    bashio::log.info "vclient_pub run finished. Looping vclient. Sleeping $REFRESH_RATE seconds."
done
