# vcontrold/rootfs/etc/services.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Fetch MQTT settings
. /etc/services.d/get_mqtt_settings.sh

## Fetch Vcontrold settings
. /etc/services.d/get_vcontrold_settings.sh

## Run your program
bashio::log.info "Starting vclient_pub main loop"
while sleep $REFRESH_RATE; do
    
#wether to run each comand separatley and open a new serial connection for each command or not
if bashio::config.true 'loop_vclient'; then 
    bashio::log.info "Starting one vclient per command"
    while read p; do
      echo "$p"
      ## only one command per run, no list anymore
        cp -f /etc/vcontrold/2_mqtt.tmpl /config/vcontrold/2_mqtt.tmpl        
        vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -c $p -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh && exitcode=0 || exitcode=$?
        cp -f /etc/vcontrold/3_mqtt_pub.sh /config/vcontrold/3_mqtt_pub.sh
        bashio::log.info "vclient_pub line $p finished with exit code $exitcode."
        sleep 2
    done < /etc/vcontrold/1_mqtt_commands.txt

else
    ## only one vclient per run, all commands are in the list 1_mqtt_commands.txt (faster, but maybe not so reliable)
    bashio::log.info "Starting one vclient for all commands"
    vclient -h $VCONTROL_HOST -p $VCONTROL_PORT -f /etc/vcontrold/1_mqtt_commands.txt -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh && exitcode=0 || exitcode=$?
fi
    bashio::log.info "vclient_pub run finished. Looping vclient. Sleeping $REFRESH_RATE seconds."
done
